# etc/skel/.config/autostart/wallpaper.desktop
[Desktop Entry]
Exec=wallpaper
Icon=application-x-shellscript
Name=wallpaper
Type=Application
X-KDE-AutostartScript=true